# vm-image/rootfs-overlay/etc/inittab
::sysinit:/sbin/mako-init
::respawn:/sbin/getty 38400 hvc0
::shutdown:/bin/umount -a -r
